# deploy/epstein-web.service
[Unit]
Description=Epstein Files Web (Next.js standalone)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/epstein/web
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_URL=https://epsteindata.cc
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/bin/node /root/epstein/web/.next/standalone/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
